# supervisord/conf/cherry-go-mgr-prod.ini
[program:cherry-go-mgr-prod]
command =  /bin/bash -c 'cd /data/server/backend/cherry-go-mgr-prod && ./cherry-go-mgr-prod'
autostart = true
startsecs = 5
autorestart = true
startretries = 3
user = root
redirect_stderr = true
stdout_logfile = /data/soft/supervisor/log/cherry-go-mgr-prod.log
stopasgroup=true
killasgroup=true

