# supervisord/conf/cherry-go-chat.ini
[program:cherry-go-chat]
command =  /bin/bash -c 'cd /data/server/cherry_go && ./home -conf=conf/app_chat.conf'
autostart = true
startsecs = 5
autorestart = true
startretries = 3
user = toney
redirect_stderr = true
stdout_logfile = /data/soft/supervisor/log/cherry-go-chat.log
stopasgroup=true
killasgroup=true

